# systemd/cdd-backend.service
[Unit]
Description=Change-Driven Development Backend
After=network.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}/backend
Environment="PATH={{PROJECT_DIR}}/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart={{PROJECT_DIR}}/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5s
StandardOutput=append:{{PROJECT_DIR}}/logs/backend.log
StandardError=append:{{PROJECT_DIR}}/logs/backend-error.log

[Install]
WantedBy=multi-user.target
